# templates/kafka-init.sh.j2
#!/bin/sh
# Copyright (C) 2016 Wesley Tanaka
### BEGIN INIT INFO
# Provides:          kafka
# Required-Start:
# Required-Stop:
# Should-Start: $network
# Default-Start:     2 3 4 5
# Default-Stop :     0 1 6
# Short-Description: Start Kafka broker
# Description:       Runs kafka-server-start.sh from kafka distribution
### END INIT INFO

{% for name, value in (kafka_env_variables.items() | list) %}
{{name}}='{{value}}'
export {{name}}
{% endfor %}

case $1 in
   start)
      LOG_DIR='{{kafka_log_directory}}'
      export LOG_DIR
      {{kafka_install_dir}}/bin/kafka-server-start.sh \
         {{kafka_install_dir}}/config/{%if 'file_basename' in item%}{{item.file_basename}}{%else%}{{kafka_default_file_basename_prefix}}{{item.broker_id}}{%endif%}.properties 2>&1 > \
         {{kafka_log_directory}}/{%if 'file_basename' in item%}{{item.file_basename}}{%else%}{{kafka_default_file_basename_prefix}}{{item.broker_id}}{%endif%}.out &
      echo $! > {{kafka_pid_directory}}{%if 'file_basename' in item%}{{item.file_basename}}{%else%}{{kafka_default_file_basename_prefix}}{{item.broker_id}}{%endif%}.pid;
      ;;
   stop)
      kill `cat {{kafka_pid_directory}}{%if 'file_basename' in item%}{{item.file_basename}}{%else%}{{kafka_default_file_basename_prefix}}{{item.broker_id}}{%endif%}.pid`
      ;;
   # Kafka does not currently have a graceful reload per
   # https://issues.apache.org/jira/browse/KAFKA-1229
   *)
      echo "usage: kafka {start|stop}"
      ;;
esac
exit 0
